# install/rethinkdb.conf
description "rethinkdb"
author "juuso haavisto"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])

stop on runlevel [!2345]
respawn

script
        exec rethinkdb
end script